# infrastructure/jenkins_ansible/web/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=10

[program:php-fpm]
command=/usr/sbin/php-fpm82 --nodaemonize
autostart=true
autorestart=true
priority=20

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
